# debian-bsps/aurora-830/update_optics_temp.timer
[Unit]
Description=Run Optics Temp Script every 10 seconds

[Timer]
OnBootSec=3s
OnUnitActiveSec=3s
AccuracySec=1s

[Install]
WantedBy=timers.target
