# debian-bsps/aurora-830/update_optics_temp.service
[Unit]
Description=Collect Max Transceiver Temperature for BMC

[Service]
Type=oneshot
ExecStart=/usr/sbin/update_optics_temp.sh
